# systemd/keyboard-context.service
[Unit]
Description=Running keyboard-context

[Service]
Type=exec
ExecStart=%h/bin/keyboard-context

[Install]
WantedBy=default.target
